# webrecorder/systemd/webrecorder-pywb.socket
[Socket]
ListenStream=/run/webrecorder-pywb.sock

[Install]
WantedBy=sockets.target